# vagrant/participant_12915.conf
[program:participant_12915]
command=/opt/helix/bin/start-helix-participant.sh --zkSvr localhost:2181 --cluster MYCLUSTER --host localhost --port 12915 --stateModelType MasterSlave
autostart=false
stderr_logfile=/var/log/supervisor/participant_12915.err.log
stdout_logfile=/var/log/supervisor/participant_12915.out.log
autorestart=true
priority=10

